# conf/supervisor/conf.d/celery_beat.conf
[program:celery-beat]
command=celery -A project beat -l info -S django
stopasgroup=true
killasgroup=true

