# conf/supervisor/conf.d/celery_worker.conf
[program:celery-worker]
command=celery -A project worker -l info
stopasgroup=true
killasgroup=true

